# deploy/agent-music-dash.service
[Unit]
Description=Agent Music operator dashboard (dash.xdshka.party)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/agent-music-tg-dash/current
# .env lives directly under /opt/agent-music-tg-dash (not inside a release
# dir) so it survives every deploy untouched. Needs: TELEGRAM_BOT_TOKEN,
# TELEGRAM_BOT_USERNAME, DASH_ADMIN_CHAT_IDS (or ADMIN_CHAT_IDS),
# DASH_SESSION_SECRET, DASH_PORT, DASH_PROD_DB_PATH, DASH_DEV_DB_PATH — see
# .env.example.
EnvironmentFile=/opt/agent-music-tg-dash/.env
ExecStart=/root/.bun/bin/bun run server/dashboard-server.ts
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
